# pruner/cron/cron.d/prune
0 3 * * * hluser /bin/bash /scripts/prune.sh >> /proc/1/fd/1 2>&1

