# pruner/cron/cron.d/s3_sync
*/15 * * * * /bin/bash -c '/home/hluser/scripts/s3_sync.sh > /proc/1/fd/1 2>&1'

